# crontab
* * * * * curl http://localhost/getcommand >/dev/null 2>&1
